# systemd/ouroboros.service
[Unit]
Description=Ouroboros Self-Evolving Autonomous Agent
Documentation=https://github.com/YOUR_USERNAME/ouroboros
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=thevetev
WorkingDirectory=/home/thevetev/ouroboros
ExecStart=/home/thevetev/ouroboros/.venv/bin/python -m ouroboros.cli moltbook run

# Auto-restart on failure
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ouroboros

# Resource limits (optional, for Raspberry Pi)
# Prevent OOM killer on low-memory systems
OOMScoreAdjust=-100

# Security hardening (optional)
# NoNewPrivileges=true
# PrivateTmp=true

[Install]
WantedBy=multi-user.target
